# base-files/etc/uci-defaults/09_broadcom-wifi-detect
#!/bin/sh

[ -d /sys/class/net/wlan0 ] && [ ! -f /etc/config/wireless ] && {
	/sbin/wifi config
	exit 0
}

exit 1
